# systemd/user/backup.timer
[Unit]
Description=Run backup.service everyday

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

