# systemd/user/backup.timer
[Unit]
Description=Automatic Restic backup on schedule

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

